# shadowsocks/etc/hotplug.d/iface/40-dropbear
#!/bin/sh

if [ "$INTERFACE" = "wan" ] && [ "$ACTION" = "ifup" ]; then
    /etc/init.d/dropbear restart
fi;
