# shadowsocks/etc/crontabs/root
# m h  dom mon dow   command
*/5 * * * * /usr/bin/tunnel_keepalive >> /var/log/tunnel.log 2>&1
*/30 * * * * /usr/bin/ss-reload
1 1 * * * /usr/bin/update
