# shadowsocks/etc/crontabs/root
# m h  dom mon dow   command
*/30 * * * * /usr/bin/ss-reload
10 4 * * * /usr/bin/update
