# shadowsocks/etc/crontabs/root
# m h  dom mon dow   command
*/30 * * * * /usr/bin/ss-reload
* 2 14 * * /usr/bin/chinalist
10 4 * * * /usr/bin/update
